# templates/jenkins@.service.erb
[Unit]
Description=Jenkins - %i - Instance
After=syslog.target

[Service]
ExecStart=/bin/sh -c "/bin/bash /usr/local/bin/jenkins_init.sh %i"
Type=forking
User=<%= @user %>

[Install]
WantedBy=multi-user.target
